# etc/cron.d/ikanow-aleph2
@hourly /bin/find /tmp/ -maxdepth 1 -name "aleph2_storm_jar_*.jar" -mtime -5|xargs rm -rf
