# rpm/etc/cron.d/ikanow-aleph2
@hourly /bin/find /tmp/ -maxdepth 1 -name "*.jar" -mtime -30|xargs rm -rf
